# crontab
SHELL=/bin/bash
*/5 * * * * sh /home/pi/raspberry-pi-timelapse/take-snapshot >> /home/pi/raspberry-pi-timelapse/snapshots.log 2>&1
@reboot /home/pi/raspberry-pi-timelapse/server
